# services/schloss/run
#!/bin/bash
cd ~/zb1/metamorphosis
./start_service schloss
